# uci-defaults/99-passwall2
#!/bin/sh

# Ensure the PassWall2 init script is executable
chmod +x /etc/init.d/passwall2

# Enable PassWall2 to start at boot
/etc/init.d/passwall2 enable

exit 0
